# .config/systemd/user/mailsync.service
[Unit]
Description=Start mailsync script
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=%h/.config/mutt/mailsync
