# .config/systemd/user/newsboat.service
[Unit]
Description=Reload newsboat feeds
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=newsboat -x reload
