# .config/systemd/user/newsboat.timer
[Unit]
Description=Reload newsboat feeds every 60 mins

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
